# molly-h5/rpm/post-install.sh
systemctl daemon-reload
systemctl enable molly-h5-server
nginx -s reload